# helpers/claude-ops.cron
# Claude ops cron
# Watchdog: every minute, scan tmux panes, fire notify on >=10 min stuck, write liveness.json
* * * * * root /usr/local/bin/claude-watchdog.sh >/dev/null 2>&1
# Respawn: every 5 min, respawn any archived project windows
*/5 * * * * root /usr/local/bin/claude-respawn.sh >/dev/null 2>&1
# Capacity: every 5 min, write capacity.json + notify on threshold breaches
*/5 * * * * root /usr/local/bin/capacity-check.sh >/dev/null 2>&1
# Daily reflection reminder at 23:00 UTC
0 23 * * * root /usr/local/bin/notify.sh "Daily reflection time — reflections/$(date -u +\%Y-\%m-\%d).md" "Platform reflection" default "https://github.com/csosin1/ClaudeCode/tree/main/reflections" >/dev/null 2>&1
# Chat telemetry: every 30 min, scan JSONLs for per-chat token cost + per-skill usage
*/30 * * * * root /usr/local/bin/chat-telemetry.py >/dev/null 2>&1
# Project check-in: every 30 min, verify active chats are making progress
*/30 * * * * root /usr/local/bin/project-checkin.sh >/dev/null 2>&1
# Branch audit: every Monday 09:00 UTC, surface stale/behind branches
0 9 * * 1 root /usr/local/bin/branch-audit.sh >/dev/null 2>&1
# Projects smoketest (report-only): every hour, records 200/regression state
17 * * * * root /usr/local/bin/projects-smoketest.sh report --quiet >/dev/null 2>&1
# Docs mirror: every 15 min, refresh /var/www/docs/ from /opt/site-deploy/
*/15 * * * * root /usr/local/bin/refresh-docs.sh >/dev/null 2>&1
# IP allow-list: weekly check for Anthropic CIDR publication (stub until they publish)
0 8 * * 1 root /usr/local/bin/refresh-docs-ip-allowlist.sh >/dev/null 2>&1
