# timeshare-surveillance/watcher/watcher.service.template
[Unit]
Description=Timeshare Surveillance EDGAR Watcher
After=network.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=__PROJECT_DIR__/.env
Environment=PYTHON_EXE=__VENV__/bin/python
Environment=HOME=/root
ExecStart=__VENV__/bin/python watcher/edgar_watcher.py
Restart=always
RestartSec=10
User=root
StandardOutput=append:/var/log/timeshare-surveillance/watcher.log
StandardError=append:/var/log/timeshare-surveillance/watcher.log

[Install]
WantedBy=multi-user.target
